# openrc/throttled
#!/sbin/openrc-run
name="throttled daemon"
command="/path/to/throttled/script"
command_backqround=true
pidfile=/var/throttled.pid
